# ops/systemd/bountyverdict-demand-watch.timer
[Unit]
Description=Poll public marketplaces for exact funded demand every ten minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=10min
RandomizedDelaySec=30s
Persistent=true
Unit=bountyverdict-demand-watch.service

[Install]
WantedBy=timers.target
