# ops/systemd/bountyverdict-directory-monitor.timer
[Unit]
Description=Run the audited BountyVerdict directory retrieval window every six hours

[Timer]
OnCalendar=*-*-* 05,11,17,23:15:00
RandomizedDelaySec=30s
Unit=bountyverdict-directory-monitor.service

[Install]
WantedBy=timers.target
